# supervisordexample.conf
[program:huepoll]
command=/usr/bin/python3 /home/homeassistant/.homeassistant/huepoll.py
directory=/home/homeassistant
autostart=true
autorestart=true
startretries=3
user=homeassistant
